# supervisor/conf.d/caddy.conf
[program:caddy]
command=caddy -conf="/home/pi/conf/Caddyfile"
autostart=true
autorestart=true

